# scripts/cyberwave-edge-ros.service
[Unit]
Description=Cyberwave Edge ROS2 Service
Documentation=https://cyberwave.com/docs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cyberwave
Group=cyberwave

# Working directory - adjust to your installation path
WorkingDirectory=/opt/cyberwave-edge-ros

# Environment file with credentials
EnvironmentFile=/opt/cyberwave-edge-ros/.env

# Source ROS 2 and run the bridge node
ExecStartPre=/bin/bash -c 'source /opt/ros/humble/setup.bash && source /opt/cyberwave-edge-ros/install/setup.bash'
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && source /opt/cyberwave-edge-ros/install/setup.bash && /opt/cyberwave-edge-ros/venv/bin/python -m ros2 run mqtt_bridge mqtt_bridge_node'

# Restart policy
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cyberwave-edge-ros

[Install]
WantedBy=multi-user.target
